# lib/v0.berlin-gtfs-rt.transport.rest-match.socket
[Unit]
Description=input for match.js

[Socket]
ListenFIFO=/run/berlin-gtfs-rt-server/1-match
SocketMode=0600
RemoveOnStop=false
