# lib/v0.berlin-gtfs-rt.transport.rest-serve.socket
[Unit]
Description=input for serve.js

[Socket]
ListenFIFO=/run/berlin-gtfs-rt-server/1-serve
SocketMode=0600
RemoveOnStop=false
